# alsa-master-pin.service
[Unit]
Description=Pin ALSA Master to 0 dB after audio device init
After=wireplumber.service
Requires=wireplumber.service

[Service]
Type=oneshot
# Wait for WirePlumber + UCM init to complete
ExecStartPre=/bin/sleep 3
# Override UCM init default of 60% with 0 dB (127/127)
ExecStart=/usr/bin/amixer -c 0 cset numid=17 127

[Install]
WantedBy=default.target
